# runner/launchd/com.orchestrator.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- TEMPLATE. Prefer scripts/setup-scheduler.sh, which substitutes HOME_DIR/REPO_PATH and installs
     the Full-Disk-Access app wrapper. KeepAlive=true => macOS auto-respawns the runner. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.orchestrator.runner</string>
  <!-- Launch the FDA-granted app wrapper; launchd cannot reliably read scripts under Documents. -->
  <key>ProgramArguments</key>
  <array>
    <string>__APP_DIR__/Contents/MacOS/ClaudeRunner</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ENABLE_PROACTIVE_LOOPS</key><string>true</string>
    <key>CLAUDE_ORCH_HOME</key><string>REPO_PATH/.runtime</string>
    <key>ORCH_LOG_DIR</key><string>REPO_PATH/.runtime/logs</string>
    <key>ORCH_LOGICAL_RUNNERS</key><string>true</string>
    <key>ORCH_KEEPALIVE_STAY_RESIDENT</key><string>true</string>
    <key>ORCH_KEEPALIVE_DUPLICATE_POLL_SECONDS</key><string>60</string>
    <key>ORCH_RUNNER_FLEET_TARGET</key><string>8</string>
    <key>ORCH_BATCH_DEV_RELEASE</key><string>true</string>
    <key>ORCH_CODE_MERGE_TARGET</key><string>dev</string>
    <key>ORCH_STAGING_BRANCH</key><string>orchestrator/dev</string>
    <key>ORCH_PUSH_ON_DEV_MERGE</key><string>true</string>
    <key>ORCH_PUSH_ON_MERGE</key><string>false</string>
    <key>ORCH_PUSH_ON_RELEASE</key><string>true</string>
    <key>RELEASE_MIN_BATCH</key><string>10</string>
    <key>RELEASE_INTERVAL_HOURS</key><string>6</string>
    <key>SESSION_AUTO_CONTINUE</key><string>true</string>
    <key>SESSION_FORCE_CONTINUE</key><string>true</string>
    <!-- add VERCEL_TOKEN / SUPABASE_ACCESS_TOKEN / provider keys here or rely on runner/.env -->
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>REPO_PATH/.runtime/logs/runner-launchd.log</string>
  <key>StandardErrorPath</key><string>REPO_PATH/.runtime/logs/runner-launchd.log</string>
</dict>
</plist>
