# runner/launchd/com.claudeorchestrator.modelscout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.claudeorchestrator.modelscout</string>
  <key>ProgramArguments</key><array>
    <string>/Applications/ClaudeRunner.app/Contents/MacOS/ClaudeRunner</string>
    <string>model_scout.py</string></array>
  <key>EnvironmentVariables</key><dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string></dict>
  <key>WorkingDirectory</key><string>/Users/kpasch/Documents/beethoven/claude-orchestrator/runner</string>
  <key>StartInterval</key><integer>14400</integer>
  <key>StandardOutPath</key><string>/Users/kpasch/Library/Logs/claude-orchestrator/modelscout.log</string>
  <key>StandardErrorPath</key><string>/Users/kpasch/Library/Logs/claude-orchestrator/modelscout.err</string>
</dict></plist>
