# scripts/launchd/com.claudeorchestrator.self-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- TEMPLATE — run scripts/setup-scheduler.sh to install. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.self-review</string>
    <key>ProgramArguments</key>
    <array>
        <string>HOME_DIR/Applications/ClaudeRunner.app/Contents/MacOS/run</string>
        <string>self_review.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SUPABASE_URL</key>
        <string>SUPABASE_URL_PLACEHOLDER</string>
        <key>SUPABASE_SERVICE_KEY</key>
        <string>SUPABASE_SERVICE_KEY_PLACEHOLDER</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>REPO_PATH/runner</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>LOG_DIR/self-review.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/self-review.err</string>
</dict>
</plist>
