# scripts/launchd/com.claudeorchestrator.deploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Canary-gated nightly deploy — runs nightly at 02:30. Requires METRICS_URL.
     TEMPLATE — run scripts/setup-scheduler.sh to install with real env vars. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.deploy</string>
    <key>ProgramArguments</key>
    <array>
        <string>__APP_DIR__/Contents/MacOS/ClaudeRunner</string>
        <string>deploy</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SUPABASE_URL</key>
        <string>SUPABASE_URL_PLACEHOLDER</string>
        <key>SUPABASE_SERVICE_KEY</key>
        <string>SUPABASE_SERVICE_KEY_PLACEHOLDER</string>
        <key>METRICS_URL</key>
        <string>METRICS_URL_PLACEHOLDER</string>
        <key>CANARY_MAX_ERROR_RATE</key>
        <string>1.0</string>
        <key>CANARY_MAX_P95_MS</key>
        <string>500</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>LOG_DIR/deploy.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/deploy.err</string>
</dict>
</plist>
