# scripts/launchd/com.claudeorchestrator.batch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Batch API off-peak pass — submits kind=batch tasks to Anthropic Batch API nightly
     at 23:30 (50% cost, ~24h turnaround). Also polls at 08:00 for previous night's results.
     Requires ANTHROPIC_API_KEY. TEMPLATE — run scripts/setup-scheduler.sh to install. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.batch</string>
    <key>ProgramArguments</key>
    <array>
        <string>__APP_DIR__/Contents/MacOS/ClaudeRunner</string>
        <string>batch</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SUPABASE_URL</key>
        <string>SUPABASE_URL_PLACEHOLDER</string>
        <key>SUPABASE_SERVICE_KEY</key>
        <string>SUPABASE_SERVICE_KEY_PLACEHOLDER</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>ANTHROPIC_API_KEY_PLACEHOLDER</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>LOG_DIR/batch.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/batch.err</string>
</dict>
</plist>
