# scripts/launchd/com.claudeorchestrator.spec.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Spec drift check — runs weekly on Sunday at 02:00.
     TEMPLATE — run scripts/setup-scheduler.sh to install with real env vars. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.spec</string>
    <key>ProgramArguments</key>
    <array>
        <string>__APP_DIR__/Contents/MacOS/ClaudeRunner</string>
        <string>spec</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SUPABASE_URL</key>
        <string>SUPABASE_URL_PLACEHOLDER</string>
        <key>SUPABASE_SERVICE_KEY</key>
        <string>SUPABASE_SERVICE_KEY_PLACEHOLDER</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>LOG_DIR/spec.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/spec.err</string>
</dict>
</plist>
