# scripts/launchd/com.claudeorchestrator.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.runner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>set -a; . /Users/kpasch/Documents/beethoven/claude-orchestrator/runner/.env; set +a; exec python3 /Users/kpasch/Documents/beethoven/claude-orchestrator/runner/runner.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kpasch/Documents/beethoven/claude-orchestrator/runner</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/kpasch/Library/Logs/claude-orchestrator/runner.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kpasch/Library/Logs/claude-orchestrator/runner.err</string>
</dict>
</plist>
