# scripts/launchd/com.claudeorchestrator.research-window.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- TEMPLATE — run scripts/setup-scheduler.sh to install with real env vars.
     Runs 2-5 AM with TASK_KIND_FILTER=research so the runner only picks up
     kind=research tasks during this window. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.research-window</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>exec timeout 10800 /usr/bin/python3 REPO_PATH/runner/runner.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SUPABASE_URL</key>
        <string>SUPABASE_URL_PLACEHOLDER</string>
        <key>SUPABASE_SERVICE_KEY</key>
        <string>SUPABASE_SERVICE_KEY_PLACEHOLDER</string>
        <key>TASK_KIND_FILTER</key>
        <string>research</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>LOG_DIR/research-window.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/research-window.err</string>
</dict>
</plist>
