# scripts/launchd/com.claudeorchestrator.self-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.self-review</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>set -a; . /Users/kpasch/Documents/beethoven/claude-orchestrator/runner/.env; set +a; python3 /Users/kpasch/Documents/beethoven/claude-orchestrator/runner/self_review.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kpasch/Documents/beethoven/claude-orchestrator/runner</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kpasch/Library/Logs/claude-orchestrator/self-review.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kpasch/Library/Logs/claude-orchestrator/self-review.err</string>
</dict>
</plist>
