# scripts/launchd/com.claudeorchestrator.agentmarket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Agent market mesh — role-aware cross-app model/vendor competition.
     TEMPLATE — run scripts/setup-scheduler.sh to install with real env vars. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudeorchestrator.agentmarket</string>
    <key>ProgramArguments</key>
    <array>
        <string>__APP_DIR__/Contents/MacOS/ClaudeRunner</string>
        <string>agentmarket</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>ORCH_SUPABASE_TIMEOUT</key>
        <string>12</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>LOG_DIR/agentmarket.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_DIR/agentmarket.err</string>
</dict>
</plist>
